# containers/nginx/supervisord-nginx.conf
[supervisord]
nodaemon=true

[program:serf]
command=serf agent -tag role=web
numprocs=1
autostart=true
autorestart=true

[program:nginx]
command=nginx
numprocs=1
autostart=true
autorestart=true
