# containers/serf/supervisord-serf.conf
[supervisord]
nodaemon=true

[program:serf]
command=serf agent -tag role=serf-agent
numprocs=1
autostart=true
autorestart=true
