# containers/php/supervisord-php-fpm.conf
[program:serf]
command=/start-serf.sh
numprocs=1
autostart=true
autorestart=true

[program:php-fpm]
command=/start-php-fpm.sh
numprocs=1
autostart=true
autorestart=true